# packaging/osx/scripts/postinstall
#!/bin/sh

set -e
/bin/launchctl load "/Library/LaunchDaemons/org.freelan.freelan.plist"
